# packages/patch-systemd-units/test/test-systemd-units.service
[Unit]
Description=This is a test
ExecStart=test-systemd-units
ExecStopPost=echo hey
[Install]
WantedBy=multi-user.target